# conf/supervisord.busy.conf
[supervisord]
nodaemon=true

[program:yo]
command = /bin/cat "YO!" > /sometext
user = root
autostart = true
